# environments/openscap/ansible/roles/openscap_scanner/templates/openscap-scan.timer.j2
[Unit]
Description=Weekly OpenSCAP compliance scan

[Timer]
OnCalendar={{ scap_schedule_weekday }} *-*-* {{ scap_schedule_hour }}:00:00
Persistent=true

[Install]
WantedBy=timers.target
