# environments/osquery-telegraf/ansible/roles/osquery_telegraf/files/packs/security_monitoring.conf
{
  "queries": {
    "installed_packages": {
      "query": "SELECT name, version, source FROM deb_packages;",
      "interval": 43200,
      "description": "Список установленных пакетов (для Debian/Ubuntu). Для CentOS/RHEL используйте 'rpm_packages'.",
      "platform": "linux"
    },
    "listening_ports": {
      "query": "SELECT p.name, p.pid, lp.port, lp.protocol, lp.address FROM listening_ports lp LEFT JOIN processes p ON lp.pid = p.pid WHERE lp.address != '127.0.0.1' AND lp.address != '::1';",
      "interval": 3600,
      "description": "Сетевые порты, слушающие на внешних интерфейсах (раз в час)."
    },
    "local_admins": {
      "query": "SELECT u.username, g.groupname FROM users u JOIN user_groups ug ON u.uid = ug.uid JOIN groups g ON ug.gid = g.gid WHERE g.groupname IN ('sudo', 'admin', 'wheel', 'root');",
      "interval": 7200,
      "description": "Пользователи в привилегированных группах (раз в 2 часа)."
    },
    "shell_history": {
      "query": "SELECT time, command FROM shell_history WHERE time > (strftime('%s', 'now') - 900);",
      "interval": 900,
      "description": "История командных оболочек за последние 15 минут.",
      "snapshot": true
    },
    "crontab": {
      "query": "SELECT command, path FROM crontab;",
      "interval": 3600,
      "description": "Задачи в системном и пользовательских crontab."
    }
  }
}
